# find_my_history_addon/rootfs/etc/services.d/find_my_history/finish
#!/usr/bin/with-contenv sh
# ==============================================================================
# Take down the S6 supervision tree when service exits
# ==============================================================================

echo "[find_my_history] Service stopped with exit code: $1"

# If the service exits with a non-zero code, take down the whole container
if [ "$1" != "0" ] && [ "$1" != "256" ]; then
    echo "[find_my_history] Service crashed, stopping container..."
    /run/s6/basedir/bin/halt
fi
